# etc/systemd/system/language_processing.service
[Unit]
Description=Gunicorn instance to serve language_processing
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/data/current
Environment=GOOD_NEWS_API_HOST=127.0.0.1
Environment=GOOD_NEWS_API_SCHEMA=http
EnvironmentFile=/data/current/credentials/goodnews-auth-credentials.env
EnvironmentFile=/data/current/credentials/goodnews-ibm-credentials.env
ExecStart=/usr/bin/env python3 -m language_processing.main
Restart=on-failure

[Install]
WantedBy=multi-user.target
